# nightly-build.cron
#!/bin/bash

# full path of the nightly build script
script_path=$(readlink -f build.sh)

# log directory
log_dir=/usr/common/contrib/bccp/build-logs

# do the nightly build (only add cron job if not present)
mailto="MAILTO=\"\""
(echo "$mailto"; ) | crontab -

for PYTHON in "2.7" "3.5" "3.6"; do
    filename="${log_dir}/conda-build-${PYTHON}-\`date +\%Y-\%m-\%d\`-cron.log"
    job="00 00 * * * bash -l $script_path $PYTHON > $filename 2>&1"
    ( crontab -l; echo "${job}"; ) | crontab -
done
